# deb/debian/postrm
#!/bin/bash -e

PERCONA_PREFIX=/usr/local
SERVICE_NAME=percona-agent
BASEDIR="${PERCONA_PREFIX}/percona/${SERVICE_NAME}"

case "$1" in
	remove)
		# delete non config files which are not automatically removed
		rm -rf ${BASEDIR}/data
		rm -rf ${BASEDIR}/trash
		rm -f ${BASEDIR}/*.log
	;;

	purge)
		# delete all files including config
		rm -rf ${BASEDIR}
	;;

	upgrade)
	;;

	*)
		echo "postrm called with unknown argument \`$1'" >&2
		exit 1
	;;
esac

# remove symlinks for service
if which update-rc.d >/dev/null 2>&1; then
 	update-rc.d percona-agent remove >/dev/null
elif which chkconfig >/dev/null 2>&1; then
    chkconfig --del percona-agent >/dev/null
fi

# dh_installdeb will replace this with shell code automatically
# generated by other debhelper scripts.

#DEBHELPER#

exit 0